# homelab/backups/templates/restic-backup.timer.j2
[Unit]
Description=Restic backup timer ({{ item.name }})

[Timer]
OnCalendar={{ backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
